# lib/systemd/system/uuidd.socket
[Unit]
Description=UUID daemon activation socket

[Socket]
ListenStream=/var/lib/run/uuidd/request

[Install]
WantedBy=sockets.target
